# conf/etc/init/bidder.conf
start on runlevel [2345]
stop on runlevel [016]

setuid bidder
setgid bidder

# chdir /home/bidder/bidder
exec authbind java -jar /home/bidder/bidder/build/libs/bidder-fat.jar